# infra/systemd/waifu-bot-scheduler.service
[Unit]
Description=Waifu Bot background tick scheduler (Dramatiq enqueue)
After=network.target redis.service waifu-bot-worker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/waifu-bot-REBORN
Environment=PYTHONPATH=/opt/waifu-bot-REBORN/src
EnvironmentFile=-/opt/waifu-bot-REBORN/.env
Environment=BACKGROUND_MODE=worker
ExecStart=/opt/waifu-bot-REBORN/.venv/bin/python -m waifu_bot.worker.scheduler
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
